# infrastructure/foss-dashboard.service
[Unit]
Description=FOSS Data Platform Dashboard
After=network.target
Wants=network.target
Requires=network.target

[Service]
Type=simple
User=gunnar
Group=gunnar
WorkingDirectory=/home/gunnar/git/foss-dataplatform
ExecStart=/home/gunnar/git/foss-dataplatform/scripts/start_dashboard.sh
ExecStop=/bin/kill -TERM $MAINPID
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy - aggressive restart on failure
Restart=always
RestartSec=3
StartLimitIntervalSec=60
StartLimitBurst=5

# Environment
Environment=FLASK_ENV=production
Environment=FLASK_DEBUG=false
Environment=PYTHONPATH=/home/gunnar/git/foss-dataplatform

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/gunnar/git/foss-dataplatform

# Resource limits
MemoryLimit=1G
CPUQuota=50%
TasksMax=100

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=foss-dashboard

[Install]
WantedBy=multi-user.target
